# spotify-yt-sync/rootfs/etc/services.d/sync/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

# One-shot sync - runs once then exits

# Read config
YOUTUBE_PLAYLIST_ID=$(bashio::config 'youtube_playlist_id')
YOUTUBE_REFRESH_TOKEN=$(bashio::config 'youtube_refresh_token')
GOOGLE_CLIENT_ID=$(bashio::config 'google_client_id')
GOOGLE_CLIENT_SECRET=$(bashio::config 'google_client_secret')
LOG_LEVEL=$(bashio::config 'log_level')

# Export for Python
export YOUTUBE_PLAYLIST_ID
export YOUTUBE_REFRESH_TOKEN
export GOOGLE_CLIENT_ID
export GOOGLE_CLIENT_SECRET
export LOG_LEVEL

cd /opt/spotify_yt_sync

bashio::log.info "Starting Spotify YouTube Sync..."

# Run sync
python3 sync.py
EXIT_CODE=$?

if [ $EXIT_CODE -eq 0 ]; then
    bashio::log.info "Sync completed successfully"
else
    bashio::log.error "Sync failed with exit code $EXIT_CODE"
fi

bashio::log.info "Add-on finished, stopping container"

# Signal s6 to stop the container
/run/s6/basedir/bin/halt
