# spotify-yt-sync/rootfs/etc/services.d/sync/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

CONFIG_PATH=/data/options.json

# Read config
YOUTUBE_PLAYLIST_ID=$(bashio::config 'youtube_playlist_id')
YOUTUBE_REFRESH_TOKEN=$(bashio::config 'youtube_refresh_token')
GOOGLE_CLIENT_ID=$(bashio::config 'google_client_id')
GOOGLE_CLIENT_SECRET=$(bashio::config 'google_client_secret')
SYNC_SCHEDULE=$(bashio::config 'sync_schedule')
LOG_LEVEL=$(bashio::config 'log_level')

# Export for Python
export YOUTUBE_PLAYLIST_ID
export YOUTUBE_REFRESH_TOKEN
export GOOGLE_CLIENT_ID
export GOOGLE_CLIENT_SECRET
export LOG_LEVEL

cd /opt/spotify_yt_sync

bashio::log.info "Spotify YouTube Sync starting..."
bashio::log.info "Schedule: ${SYNC_SCHEDULE}"

# Run initial sync
bashio::log.info "Running initial sync..."
python3 sync.py

# Setup cron schedule
echo "${SYNC_SCHEDULE} cd /opt/spotify_yt_sync && python3 sync.py >> /proc/1/fd/1 2>&1" > /etc/crontabs/root

# Start crond in foreground
bashio::log.info "Starting scheduler..."
exec crond -f -l 2
